# build/qbittorrent.conf
[program:watchdog-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/qbittorrent.sh
umask = 000


